# rurouni-cache.init
#!/usr/bin/env bash
# chkconfig: 345 95 20
# description: rurouni-cache daemon
# processname: rurouni-cache

source /opt/kenshin/venv/bin/activate
mount_num=6

for instance in $(seq 0 $[(`nproc`-1)*$mount_num])
do
    eval '/opt/kenshin/bin/rurouni-cache.py --config=/etc/kenshin/rurouni.conf --instance=${instance} $@'
done
exit $?

